# infra/local/com.ourosss.hermes-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.ourosss.hermes-sync</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/mkurkin/experiments/ourosss/infra/local/sync-push.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/mkurkin/Library/Logs/ourosss-sync.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mkurkin/Library/Logs/ourosss-sync.err.log</string>
  </dict>
</plist>
